# turbo-crl.cron
22 */4 * * *	root	/usr/bin/turbo-crl.py -d 1200 /etc/grid-security/certificates
